# modules/solr/templates/solr.conf
[program:solr]
command=java -jar /opt/solr/default/start.jar
directory=/opt/solr/default
user=root
group=root
stopsignal=TERM
