# modules/sentry/templates/supervisor.udp.conf.erb
[program:sentry-udp]
command=<%= @sentry_ve_dir %>/bin/sentry --config=/etc/sentry.conf.py start udp
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=TERM
