# modules/sentry/templates/supervisor.udp.conf.erb
[program:sentry-udp]
command=sentry --config=/etc/sentry.conf.py start udp
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=TERM
