# modules/sentry/templates/upstart.udp.conf.erb
description "Sentry (udp)"
start on runlevel [2345]
stop on runlevel [06]
exec sentry --config=/etc/sentry.conf.py start udp
