# modules/openresty/files/openresty.conf
description "OpenResty"
author "Arcus, Inc."

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec /usr/local/openresty/nginx/sbin/nginx -c /usr/local/openresty/nginx/conf/nginx.conf
