# modules/graylog2/templates/graylog.conf.erb
[program:graylog-web]
command=/usr/bin/ruby script/rails server
    -e production -p <%= @graylog_web_port %>
directory=<%= scope.lookupvar('graylog2::params::graylog_dir') %>
user=root
stopsignal=QUIT
